# supervisor/container.ini
[supervisord]
nodaemon=true

[program:finch]
command=finch

